# ops/launchd/com.agent-infra.codebase-map-refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent-infra.codebase-map-refresh</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/alien/Projects/agent-infra/scripts/refresh-codebase-maps.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/alien/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/alien</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/alien/Projects/agent-infra</string>
    <!-- 06:30 — after friend-sync's 06:00 repo pull, so maps reflect fresh code. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/alien/.claude/logs/codebase-map-refresh.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/alien/.claude/logs/codebase-map-refresh.err</string>
    <key>LowPriorityIO</key>
    <true/>
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
